# .platform/wsgi.d/00_application.conf
[program:application]
command=gunicorn app:application --bind 127.0.0.1:8000
user=webapp
autostart=true
autorestart=true
redirect_stderr=true
